# cp/conf/cron.d/lepuscp
01 12 * * * root certbot renew --quiet --allow-subset-of-names
*/5 * * * *  root /usr/local/lepuscp/start.bash >/dev/null 2>&1
*/10 * * * * root php /usr/local/lepuscp/php/ssl.php >/dev/null 2>&1
0 12 * * * root find /usr/local/lepuscp/sess/* -type f -mtime +7 -delete >/dev/null 2>&1
